# ha-mcp/rootfs/etc/services.d/ha-mcp/finish
#!/usr/bin/execlineb -S0
# Arrêt propre du container quand le service s'arrête
