# ha-mcp/rootfs/etc/services.d/ha-mcp/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# Délègue au script principal qui charge les env vars depuis bashio::config
exec /run.sh
